# mkosi.profiles/ironcore/mkosi.extra/usr/lib/systemd/system/ironcore-hostname.service
[Unit]
Description=Set hostname from IronCore metadata
After=ironcore-metadata.service
Requires=ironcore-metadata.service
Before=ironcore-bootstrap.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/ironcore/ironcore-hostname.sh

[Install]
WantedBy=ironcore-bootstrap.target
