# mkosi.profiles/ironcore/mkosi.extra/usr/lib/systemd/system/ironcore-bootstrap.target
[Unit]
Description=IronCore Bootstrap
ConditionFirstBoot=yes
ConditionPathExists=!/etc/initrd-release

[Install]
WantedBy=multi-user.target
